# chromatik/pi_setup/brightness-toggle.service
[Unit]
Description=Entwined Brightness Toggle
After=network.target

[Service]
ExecStart=/home/ent/Entwined/chromatik/brightness-toggle.sh
WorkingDirectory=/home/ent/Entwined/chromatik
StandardOutput=inherit
StandardError=inherit
Restart=always
User=root

[Install]
WantedBy=multi-user.target
